# deploy/systemd/revops-poll.timer
[Unit]
Description=RevOps Command Center poll — every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
