# deploy/systemd/revops-digest.service
[Unit]
Description=RevOps Command Center — weekly Slack deal-health digest
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/revops-command-center
ExecStart=/root/revops-command-center/bin/post-digest.sh
